# ansible/roles/wg-manager/templates/wg-manager.service.j2
# {{ ansible_managed }}
[Unit]
Description=WireGuard Manager Service
After=network.target wg-quick@{{ wg_manager_wg_interface }}.service
Wants=wg-quick@{{ wg_manager_wg_interface }}.service

[Service]
Type=simple
User={{ wg_manager_user }}
Group={{ wg_manager_group }}
WorkingDirectory={{ wg_manager_install_dir }}
Environment=PATH={{ wg_manager_venv_path }}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart={{ wg_manager_venv_path }}/bin/python {{ wg_manager_install_dir }}/app.py
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
